# ops/launchd/com.lban.ai-content.x-following-filter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.lban.ai-content.x-following-filter</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/lban/.nvm/versions/node/v22.14.0/bin/node</string>
      <string>/Users/lban/dev/docs/ai-content/scripts/run-x-following-filter.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/lban/dev/docs/ai-content</string>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/lban/dev/docs/ai-content/.runtime/x-following-filter/logs/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/lban/dev/docs/ai-content/.runtime/x-following-filter/logs/launchd.stderr.log</string>
  </dict>
</plist>
